# frigate_proxy/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs the NGINX daemon
# ==============================================================================

bashio::log.info "Starting NGINX..."
exec nginx
