# frigate_proxy/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Configures NGINX for use with this add-on.
# ==============================================================================

# Get the ingress entry path
ingress_entry=$(bashio::addon.ingress_entry)

# Merge the ingress entry into the full configuration and process the template
# We use jq to add the 'entry' field to the existing options.json
jq --arg entry "${ingress_entry}" \
    '. + {entry: $entry}' /data/options.json \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf
