# frigate_proxy/rootfs/etc/s6-overlay/s6-rc.d/init-landing/run
#!/command/with-contenv bashio
# ==============================================================================
# Process landing page template
# ==============================================================================

# Process the landing page template
tempio \
    -conf /data/options.json \
    -template /etc/nginx/templates/html/landing.html.gtpl \
    -out /etc/nginx/html/landing.html
