# frigate_proxy/rootfs/etc/s6-overlay/s6-rc.d/init-landing/up
#!/command/execlineb -P
/etc/s6-overlay/s6-rc.d/init-landing/run
